# bluelog.conf
[program:bluelog]
command=pipenv shell  gunicorn -w 4 wsgi:app
directory=/home/bluelog 
user=root
autorestart=true 
autostart=true
stopasgroup=true
killasgroup=true 